# ops/systemd/velmix-queue-restart.service
[Unit]
Description=VELMiX queue restart hook
After=network.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/velmix/current
EnvironmentFile=-/etc/velmix/velmix.env
Environment=APP_ENV=production
ExecStart=/usr/bin/php artisan queue:restart
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/www/velmix/current/storage /var/www/velmix/current/bootstrap/cache /var/www/velmix/shared
UMask=0027
